# deb/usr/lib/systemd/user/power-mode-auto-changer.service
[Unit]
Description=Sync power mode with current power source

[Service]
Type=oneshot
ExecStart=/usr/local/bin/power-mode-auto-changer-service

[Install]
WantedBy=default.target