# deb/DEBIAN/prerm
#!/bin/bash
set -e

systemctl stop power-mode-auto-changer.service || true
systemctl disable power-mode-auto-changer.service || true

exit 0