# deb/lib/udev/rules.d/99-power-mode-auto-changer.rules
SUBSYSTEM=="power_supply", ATTR{type}=="Mains", RUN+="/usr/local/lib/power-mode-auto-changer/switch-power-mode.sh"